# .config/systemd/user/apps.target
[Unit]
Description=Apps to start on login
Wants=wm.target
AllowIsolate=true

[Install]
Alias=default.target
